# examples/scenarios/secret_apt29_emulation.json
{
  "scenario_id": "scen-secret-apt29-001",
  "version": "1.0.0",
  "name": "APT29 Cozy Bear Full Kill Chain",
  "description": "Complete APT29 adversary emulation with credential theft, lateral movement, and data collection",
  "classification_level": "SECRET",
  "author": "red-team-classified",
  "created_at": "2025-01-14T10:00:00Z",
  "persona_reference": {
    "persona_id": "550e8400-e29b-41d4-a716-446655440100",
    "name": "APT29 Cozy Bear",
    "attribution": "Russia/SVR"
  },
  "targets": [
    {
      "name": "enclave-dc-01",
      "type": "domain_controller",
      "network": "192.168.10.0/24",
      "ip": "192.168.10.10",
      "criticality": "critical",
      "services": ["ldap", "kerberos", "dns"]
    },
    {
      "name": "enclave-srv-file-01",
      "type": "file_server",
      "network": "192.168.10.0/24",
      "ip": "192.168.10.20",
      "criticality": "high",
      "services": ["smb", "rdp"]
    },
    {
      "name": "enclave-srv-mail-01",
      "type": "mail_server",
      "network": "192.168.10.0/24",
      "ip": "192.168.10.30",
      "criticality": "high",
      "services": ["exchange"]
    },
    {
      "name": "sipr-ws-001",
      "type": "workstation",
      "network": "192.168.20.0/24",
      "ip": "192.168.20.100",
      "criticality": "medium",
      "owner": "target-user-001"
    },
    {
      "name": "sipr-ws-002",
      "type": "workstation",
      "network": "192.168.20.0/24",
      "ip": "192.168.20.101",
      "criticality": "medium",
      "owner": "target-user-002"
    }
  ],
  "kill_chain": [
    "reconnaissance",
    "initial_access",
    "execution",
    "persistence",
    "privilege_escalation",
    "defense_evasion",
    "credential_access",
    "discovery",
    "lateral_movement",
    "collection"
  ],
  "phases": [
    {
      "phase_id": "phase-recon-001",
      "name": "Internal Reconnaissance",
      "type": "reconnaissance",
      "objective": "Map internal network and identify high-value targets",
      "actions": [
        {
          "action_id": "action-001",
          "type": "discovery",
          "technique_id": "T1087",
          "tool": "bloodhound",
          "parameters": {
            "collection_method": "all",
            "stealth": true
          }
        },
        {
          "action_id": "action-002",
          "type": "discovery",
          "technique_id": "T1482",
          "tool": "adfind",
          "parameters": {
            "query": "trusts,users,computers"
          }
        }
      ]
    },
    {
      "phase_id": "phase-privesc-001",
      "name": "Privilege Escalation",
      "type": "privilege_escalation",
      "objective": "Obtain domain admin privileges",
      "depends_on": ["phase-recon-001"],
      "actions": [
        {
          "action_id": "action-003",
          "type": "credential_access",
          "technique_id": "T1558",
          "tool": "rubeus",
          "parameters": {
            "attack": "kerberoast",
            "output_format": "hashcat"
          }
        },
        {
          "action_id": "action-004",
          "type": "credential_access",
          "technique_id": "T1003",
          "tool": "mimikatz",
          "parameters": {
            "command": "sekurlsa::logonpasswords"
          }
        }
      ]
    },
    {
      "phase_id": "phase-lateral-001",
      "name": "Lateral Movement",
      "type": "lateral_movement",
      "objective": "Move to file server and mail server",
      "depends_on": ["phase-privesc-001"],
      "actions": [
        {
          "action_id": "action-005",
          "type": "lateral_movement",
          "technique_id": "T1021",
          "tool": "impacket",
          "parameters": {
            "method": "wmiexec",
            "target": "enclave-srv-file-01"
          }
        },
        {
          "action_id": "action-006",
          "type": "lateral_movement",
          "technique_id": "T1550",
          "tool": "pth_toolkit",
          "parameters": {
            "method": "pass_the_hash",
            "target": "enclave-srv-mail-01"
          }
        }
      ]
    },
    {
      "phase_id": "phase-persist-001",
      "name": "Persistence",
      "type": "persistence",
      "objective": "Establish persistent access",
      "depends_on": ["phase-lateral-001"],
      "actions": [
        {
          "action_id": "action-007",
          "type": "persistence",
          "technique_id": "T1053",
          "tool": "schtasks",
          "parameters": {
            "task_name": "SystemHealthCheck",
            "trigger": "daily",
            "time": "03:00"
          }
        }
      ]
    },
    {
      "phase_id": "phase-c2-001",
      "name": "C2 Establishment",
      "type": "defense_evasion",
      "objective": "Establish covert command and control",
      "depends_on": ["phase-persist-001"],
      "actions": [
        {
          "action_id": "action-008",
          "type": "command_and_control",
          "technique_id": "T1071",
          "tool": "cobalt_strike",
          "parameters": {
            "beacon_type": "https",
            "sleep": 60,
            "jitter": 20
          }
        }
      ]
    },
    {
      "phase_id": "phase-collect-001",
      "name": "Data Collection",
      "type": "collection",
      "objective": "Collect target documents and emails",
      "depends_on": ["phase-lateral-001"],
      "actions": [
        {
          "action_id": "action-009",
          "type": "collection",
          "technique_id": "T1213",
          "tool": "snaffler",
          "parameters": {
            "paths": ["\\\\enclave-srv-file-01\\shares"],
            "keywords": ["classified", "secret", "confidential"]
          }
        },
        {
          "action_id": "action-010",
          "type": "collection",
          "technique_id": "T1114",
          "tool": "mailsniper",
          "parameters": {
            "search": "subject:secret OR body:classified"
          }
        }
      ]
    }
  ],
  "ttps": {
    "primary_techniques": ["T1558", "T1003", "T1021", "T1213"],
    "evasion_techniques": ["T1027", "T1070", "T1036"],
    "tools": ["bloodhound", "mimikatz", "rubeus", "cobalt_strike"]
  },
  "success_criteria": {
    "domain_admin_obtained": true,
    "lateral_movement_complete": true,
    "data_collected": true,
    "c2_established": true,
    "persistence_achieved": true
  },
  "deconfliction": {
    "expected_alerts": [
      "Kerberoasting detected",
      "DCSync attempt",
      "Suspicious scheduled task creation"
    ],
    "soc_contact": "soc-classified@example.gov",
    "abort_conditions": [
      "Production system impact",
      "Detection by external party"
    ]
  },
  "metadata": {
    "estimated_duration_hours": 72,
    "stealth_level": "medium",
    "complexity": "high",
    "tags": ["apt29", "full-kill-chain", "credential-theft", "lateral-movement"]
  }
}
